# deploy/evenhub/systemd/nanoclaw.service.d/evenhub.conf
[Unit]
Wants=network-online.target
Requires=nanoclaw-moonshine.service
After=network-online.target nanoclaw-moonshine.service

[Service]
EnvironmentFile=/etc/nanoclaw/evenhub.env
