# deploy/evenhub/systemd/nanoclaw-whisper.service
[Unit]
Description=NanoClaw local Whisper inference
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=nanoclaw-whisper
Group=nanoclaw-whisper
ExecStart=/usr/local/bin/whisper-server --model /var/lib/nanoclaw/whisper/ggml-base.en.bin --host 127.0.0.1 --port 8178 --threads 4 --processors 1 --language en --no-timestamps --no-context
Restart=on-failure
RestartSec=2
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
RestrictSUIDSGID=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
IPAddressDeny=any
IPAddressAllow=localhost
ReadOnlyPaths=/var/lib/nanoclaw/whisper
# Initial ceiling for base.en; ticket 05 replaces it with measured peak RSS + 25%.
MemoryMax=1G

[Install]
WantedBy=multi-user.target
