# deploy/evenhub/systemd/nanoclaw-evenhub-firewall.service
[Unit]
Description=NanoClaw EvenHub LAN firewall boundary
Before=caddy.service nanoclaw-tailscale-serve.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/nft -f /etc/nftables.d/nanoclaw-evenhub.nft
ExecReload=/usr/sbin/nft -f /etc/nftables.d/nanoclaw-evenhub.nft
RemainAfterExit=yes
NoNewPrivileges=true
CapabilityBoundingSet=CAP_NET_ADMIN
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
